# source/multimedia-extra/oss/rc.oss
#!/bin/bash

# (c) 2010 James Buren <ryuo@frugalware.org>
# rc.oss for Frugalware
# distributed under GPL License

# chkconfig: 2345 99 01
# description: Starts and stops OSS

source /lib/initscripts/functions
TEXTDOMAIN=oss
TEXTDOMAINDIR=/lib/initscripts/messages

actions=(restart status start stop)
daemon=$"OSS"

rc_start()
{
	start_msg
	/usr/sbin/soundon
	ok $?
}

rc_stop()
{
	stop_msg
	/usr/sbin/soundoff
	ok $?
}

rc_exec $1

# vim: ft=sh
