# source/base/splashy/rc.zz-splash
#!/bin/bash

# (c) 2006 Miklos Vajna <vmiklos@frugalware.org>
# rc.zz-splash for Frugalware
# distributed under GPL License

# chkconfig: 2345 99 01
# description: Shuts down the splash screen at the end of the boot
source /lib/initscripts/functions
actions=(start stop)

rc_start()
{
	pidof splashy >/dev/null && splashy_update exit
}

rc_stop()
{
	:
}

rc_exec $1
