# source/base/util-linux/rc.console
#!/bin/bash

# (c) 2005      Marcus Habermehl <bmh198de@yahoo.de>
# (c) 2003-2005 Vajna Miklos <vmiklos@frugalware.org>
# rc.console for Frugalware
# distributed under GPL License

# chkconfig: 2345 00 98
# description: Screen blanks after n minutes idle time.

source /lib/initscripts/functions

actions=(start)

rc_start()
{
	if [ -e /etc/sysconfig/console ] ; then
		source /etc/sysconfig/console
		setterm -blank $idle
	fi
}

rc_exec $1

# vim: ft=sh
