# source/base/util-linux/rc.console
#!/bin/sh

# (c) 2003 Vajna Miklos <mamajom@axelero.hu>
# rc.console for FrugalWare
# distributed under GPL License

#. /etc/rc.d/rc.functions

# chkconfig: 2345 00 98
# description: Screen blanks after n minutes idle time

if [ "$1" = "start" ]; then
	[ -e /etc/sysconfig/console ] || exit 1
	. /etc/sysconfig/console
	setterm -blank $idle
fi
