# source/base/util-linux/rc.console
#!/bin/sh

# (c) 2003-2005 Vajna Miklos <vmiklos@frugalware.org>
# rc.console for Frugalware
# distributed under GPL License

# chkconfig: 2345 00 98
# description: Screen blanks after n minutes idle time.

if [ "$1" = "start" ]; then
	[ -e /etc/sysconfig/console ] || exit 1
	. /etc/sysconfig/console
	setterm -blank $idle
fi
