# source/base/util-linux/rc.rmount
#!/bin/sh

# (c) 2003 Vajna Miklos <mamajom@axelero.hu>
# rc.rmount for FrugalWare
# distributed under GPL License

. /etc/rc.d/rc.functions

# chkconfig: 2345 19 90
# description: Mounts or umounts remote (NFS, SMB) filesystems

if [ "$1" = "stop" ]; then
	stop "$stoprmount"
	umount -a -r -t nfs,smbfs
	ok $?
elif [ "$1" = "restart" ]; then
	"$0" stop
	sleep 1
	"$0" start
else # start
	start "$startrmount"
	mount -a -t nfs,smbfs
	ok $?
fi
