# source/apps-extra/festival/rc.festival
#!/bin/sh

# (c) 2005 Miklos Vajna <vmiklos@frugalware.org>
# rc.acpid for Frugalware
# distributed under GPL License

. /etc/rc.d/rc.functions

# chkconfig: 2345 44 56
# description: Start/stop festival server

if [ "$1" = "stop" ]; then
	stop "$stopfestival"
	killall festival
	ok $?
elif [ "$1" = "restart" ]; then
	"$0" stop
	sleep 1
	"$0" start
else # start
	start "$startfestival"
	festival --server -b /etc/festival/server.scm &>/dev/null &
	ok $?
fi
