# source/devel/texinfo/texinfo.cron
#!/bin/sh

mkinfodir /usr/share/info > /usr/share/info/dir
