# source/network-extra/inetd/inetd.service
[Unit]
Description=Internet super-server daemon

[Service]
ExecStart=/usr/sbin/inetd
Type=forking

[Install]
WantedBy=multi-user.target
