# source/apps-extra/firebird/firebird-classic.socket
[Unit]
Description=Firebird Classic Activation Socket

[Socket]
ListenStream=3050
Accept=true

[Install]
WantedBy=sockets.target

