# source/base/kernel/cpupower.service
[Unit]
Description=Apply cpupower configuration
After=syslog.target

[Service]
ExecStart=/usr/sbin/cpupower.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
