# source/network-extra/amavisd-new/amavisd-new.service
[Unit]
Description=Amavis Mail Filter Daemon
After=network.target

[Service]
ExecStart=/usr/libexec/amavisd-wrapper
PIDFile=/var/lock/amavis/amavisd.pid
Type=forking

[Install]
WantedBy=multi-user.target
