# extra/source/network/pdns-recursor/rc.pdns-recursor
#!/bin/sh

# (c) 2005 Zoltan Kiss <djsmiley@frugalware.org>
# rc.pdns-recursor for FrugalWare
# distributed under GPL License

. /etc/rc.d/rc.functions

# chkconfig: 345 80 75
# description: Power DNS Recursor Startup

prefix=/usr
exec_prefix=${prefix}
BINARYPATH=${exec_prefix}/bin
SBINARYPATH=${exec_prefix}/sbin
PIDFILE=/var/run/pdns_recursor.pid
SOCKETPATH=/var/run

if [ -f $PIDFILE ]
then

if [ -s $PIDFILE -a -d /proc/$(cat $PIDFILE) 2>/dev/null ] 
then
	RUNNING=1
	PID=$(cat $PIDFILE)
else
	RUNNING=0
fi
else
	RUNNING=0
fi


if [ "$1" = "stop" ]; then
		stop "$stopprecursor"
		if [ $RUNNING = "1" ]
		then
			kill -9  $PID
		    ok $?
		else
		    ok 1
		fi

	elif [ "$1" = "restart" ]; then
	        "$0" stop
	        sleep 1
	        "$0" start
	else # start
		start "$startprecursor"	
		if [ $RUNNING = "1" ]
		then
		    ok 1
		else
			$SBINARYPATH/pdns_recursor --daemon > /dev/null 2> /dev/null
		    ok $?

		fi
fi
