# extra/source/network/pdns-recursor/rc.pdns-recursor
#!/bin/bash

# chkconfig: 2345 99 01 # Change me! :)
# description: <pdns-recursor  init scripts> 

source /lib/initscripts/functions
TEXTDOMAIN=pdns-recursor # Change me! :)
TEXTDOMAINDIR=/lib/initscripts/messages

actions=(start stop) # Change me! :)
daemon=$"pdns-recursor" # Change me! :)

prefix=/usr
exec_prefix=${prefix}
BINARYPATH=${exec_prefix}/bin
SBINARYPATH=${exec_prefix}/sbin
PIDFILE=/var/run/pdns_recursor.pid
SOCKETPATH=/var/run


pid="pidof pdns_recursor 2> /dev/null" # Change me! :)

# Some scripts need the new runlevel which is contained in $2
if [ -z $2 ] ; then
	rl=4 # If we doesn't have a runlevel, we use our default
else
	rl=$2
fi

# For available functions have a look on /source/base/sysvinit/functions

rc_start()
{
	start_msg
	if [[ "$(eval $pid)" == "" ]] ; then
		$SBINARYPATH/pdns_recursor --daemon > /dev/null 2> /dev/null
		ok $?
	else
		ok 999
	fi
}

rc_stop()
{
	stop_msg
	if [[ "$(eval $pid)" != "" ]] ; then
			kill -9  $pid
		ok $?
	else
		ok 999
	fi
}

rc_exec $1
