# source/network/bind/bind.service
[Unit]
Description=Domain Name Server
After=network.target

[Service]
ExecStart=/usr/sbin/named
Type=forking

[Install]
WantedBy=multi-user.target
