# source/network/bind/rc.bind
#!/bin/sh

# (c) 2003 Vajna Miklos <mamajom@axelero.hu>
# rc.bind for FrugalWare
# distributed under GPL License

. /etc/rc.d/rc.functions

# chkconfig: 2345 17 45
# description: named (BIND) is a Domain Name Server (DNS) \
#	       that is used to resolve host names to IP addresses.

if [ "$1" = "stop" ]; then
	stop "$stopbind"
	killall named
	ok $?
elif [ "$1" = "restart" ]; then
	"$0" stop
	sleep 1
	"$0" start
else # start
	start "$startbind"
	/usr/sbin/named
	ok $?
fi
