# source/network/dhcp/rc.dhcpd
#!/bin/sh

# (c) 2003-2004 Vajna Miklos <mamajom@axelero.hu>
# rc.dhcpd for Frugalware
# distributed under GPL License

. /etc/rc.d/rc.functions

# chkconfig: 2345 65 35
# description: dhcpd provide access to Dynamic Host Control Protocol.

if [ "$1" = "stop" ]; then
	stop "$stopdhcp"
	killall dhcpd
	ok $?
elif [ "$1" = "restart" ]; then
	"$0" stop
	sleep 1
	"$0" start
else # start
	start "$startdhcp"
	dhcpd
	ok $?
fi
