# source/network/inetd/rc.inetd
#!/bin/sh

# (c) 2003 Vajna Miklos <mamajom@axelero.hu>
# rc.inetd for FrugalWare
# distributed under GPL License

. /etc/rc.d/rc.functions

# chkconfig: 2345 26 50
# description: Internet super-server - ported from OpenBSD

if [ "$1" = "stop" ]; then
	stop "$stopinetd"
	killall inetd
	ok $?
elif [ "$1" = "restart" ]; then
	"$0" stop
	sleep 1
	"$0" start
else # start
	start "$startinetd"
	/usr/sbin/inetd
	ok $?
fi
