# source/base/apmd/rc.apmd
#!/bin/sh

# (c) 2003 Vajna Miklos <mamajom@axelero.hu>
# rc.apmd for FrugalWare
# distributed under GPL License

. /etc/rc.d/rc.functions

# chkconfig: 345 44 56
# description: Listen and dispatch APM events from the kernel

cat /proc/apm 1> /dev/null 2> /dev/null || exit 1

if [ "$1" = "stop" ]; then
	stop "$stopapm"
	killall apmd
	ok $?
elif [ "$1" = "restart" ]; then
	"$0" stop
	sleep 1
	"$0" start
else # start
	start "$startapm"
	apmd
	ok $?
fi
