# docs/skel/skel.service
[Unit]
Description=Skeleton service
After=syslog.target

[Service]
ExecStart=/bin/true
Type=forking

[Install]
WantedBy=multi-user.target
