# source/xapps-extra/virtualbox/vboxservice.service
[Unit]
Description=VirtualBox Guest Service
ConditionVirtualization=oracle
Wants=vboxdrmclient.path
Before=display-manager.service

[Service]
ExecStartPre=-/sbin/modprobe vboxguest
ExecStartPre=-/sbin/modprobe vboxvideo
ExecStartPre=-/sbin/modprobe vboxsf
ExecStart=/usr/bin/VBoxService -f

[Install]
WantedBy=multi-user.target
