# source/apps/powernowd/rc.powernowd
#!/bin/sh

# (c) 2005 Miklos Vajna <vmiklos@frugalware.org>
# rc.powernowd for Frugalware
# distributed under GPL License

. /etc/rc.d/rc.functions

# chkconfig: 2345 44 56
# description: A small daemon to adjust cpu speed.

[ -d /sys/devices/system/cpu/ ] || exit 1

if [ "$1" = "stop" ]; then
	stop "$stoppowernowd"
	killall powernowd
	ok $?
elif [ "$1" = "restart" ]; then
	"$0" stop
	sleep 1
	"$0" start
else # start
	start "$startpowernowd"
	powernowd -q
	ok $?
fi
