# source/devel/mono/rc.mono.old
#!/bin/sh

# (c) 2004, 2005 Laszlo Dvornik <dvornik@gnome.hu>
# rc.mono for Frugalware
# distributed under GPL License

. /etc/rc.d/rc.functions

# chkconfig: 2345 99 02
# description: Registering/unregistering .NET IL binaries with mono.

[ -d /proc/sys/fs/binfmt_misc ] || exit 1

if [ "$1" = "stop" ]; then
	stop "$stopmono"
	if [ -f /proc/sys/fs/binfmt_misc/CLR ]; then
	    echo '-1' > /proc/sys/fs/binfmt_misc/CLR
	fi
	ok $?
elif [ "$1" = "restart" ]; then
	"$0" stop
	sleep 1
	"$0" start
else # start
	start "$startmono"
	mount -t binfmt_misc binfmt_misc /proc/sys/fs/binfmt_misc
	echo ':CLR:M::MZ::/usr/bin/mono:' > /proc/sys/fs/binfmt_misc/register
	ok $?
fi
