# source/network/openntpd/rc.ntpd
#!/bin/sh

# (c) 2004 Peter BARABAS <z0d@artifact.hu>
# rc.ntpd for FrugalWare
# distributed under the GPL License

. /etc/rc.d/rc.functions

# chkconfig: 2345 32 56
# description: Starts or stops OpenNTPD

if [ "$1" = "stop" ]; then
	stop "$stopntpd"
	killall ntpd
	ok $?
elif [ "$1" = "restart" ]; then
	stop "$stopntpd"
	killall ntpd
	sleep 2
	"$0" start
else # start
	start "$startntpd"
	/usr/sbin/ntpd
	ok $?
fi
