# extra/source/lib/hpoj/rc.ptal
#! /bin/bash

. /etc/rc.d/rc.functions

# chkconfig: 2345 59 61
# description: ptal-init is the init script and device setup utility \
#              for the HP OfficeJet Linux driver.

if [[ "$1" == "stop" ]] ; then
	stop "$stopptal"
	/usr/sbin/ptal-init stop > /dev/null
	ok $?
elif [[ "$1" == "restart" ]] ; then
	"$0" stop
	sleep 1
	"$0" start
else # start
	start "$startptal"
	/usr/sbin/ptal-init start > /dev/null
	ok $?
fi
