# source/network-extra/nginx/nginx.service
[Unit]
Description=nginx daemon
After=network.target

[Service]
ExecStart=/usr/sbin/nginx
PIDFile=/run/nginx.pid
Type=forking

[Install]
WantedBy=multi-user.target
