# source/base/rp-pppoe/rc.adsl
#!/bin/sh

# (c) 2003 Vajna Miklos <vmiklos@frugalware.org>
# rc.adsl for FrugalWare
# distributed under GPL License

. /etc/rc.d/rc.functions

# chkconfig: 2345 15 85
# description: Init script to bring up a PPPoE link

if [ "$1" = "stop" ]; then
	stop "$stopadsl"
	pppoe-stop &>/dev/null
	ok $?
elif [ "$1" = "restart" ]; then
	"$0" stop
	sleep 1
	"$0" start
else # start
	start "$startadsl"
	pppoe-start &>/dev/null
	retval=$?
	ok $retval
fi
