# source/xapps/speech-dispatcher/speechd.service
[Unit]
Description=speechd
After=syslog.target

[Service]
ExecStart=/usr/bin/speech-dispatcher -d -P /run/speechd.pid
PIDFile=/run/speechd.pid
Type=forking

[Install]
WantedBy=multi-user.target
