# extra/source/network/piespy/rc.piespy
#! /bin/bash

. /etc/rc.d/rc.functions

# chkconfig: 2345 99 01
# description: Starting the PieSpy Social Network Bot

if [[ "$1" == "stop" ]] ; then
	stop "$stoppiespy"
	PID=`ps -ef | grep "java.*SocialNetworkBot" | grep -v grep | awk '{print $2}'`
	kill -HUP $PID
	ok $?
elif [[ "$1" == "restart" ]] ; then
	"$0" stop
	sleep 1
	"$0" start
else # start
	start "$startpiespy"
	cd /usr/lib/piespy && sh run.sh 2>&1 >> /var/log/piespy.log &
	ok $?
fi
