# source/base/cups/rc.cups
#!/bin/sh

# (c) 2003 Vajna Miklos <mamajom@axelero.hu>
# rc.cups for FrugalWare
# distributed under GPL License

. /etc/rc.d/rc.functions

# chkconfig: 2345 99 01
# description: Startup/shutdown script for the Common UNIX \
#              Printing System (CUPS).

if [ "$1" = "stop" ]; then
	stop "$stopcups"
	killall cupsd
	ok $?
elif [ "$1" = "restart" ]; then
	"$0" stop
	sleep 1
	"$0" start
else # start
	start "$startcups"
	/usr/sbin/cupsd
	ok $?
fi
