# source/devel/scratchbox/rc.scratchbox
#!/bin/bash

# (c) 2006 Miklos Vajna <vmiklos@frugalware.org>
# rc.scratchbox for Frugalware
# distributed under GPL License

# chkconfig: 2345 99 01
# description: Start/stop scratchbox

source /lib/initscripts/functions

actions=(restart start stop)
daemon="Scratchbox"

rc_start()
{
	start_msg
	/usr/lib/scratchbox/sbin/sbox_ctl start >/dev/null
	ok $?
}

rc_stop()
{
	stop_msg
	/usr/lib/scratchbox/sbin/sbox_ctl stop >/dev/null
	ok $?
}

rc_exec $1
