# source/network-extra/netkit-rwho/rwhod.service
[Unit]
Description=rwho daemon
After=network.target

[Service]
ExecStart=/usr/sbin/rwhod
Type=forking

[Install]
WantedBy=multi-user.target
