# source/base/inetutils/rlogin@.service
[Unit]
Description=Remote Login Facilities Server
After=local-fs.target

[Service]
ExecStart=-/usr/lib/inetutils/rlogind
StandardInput=socket
