# source/base/inetutils/telnet@.service
[Unit]
Description=Telnet Server
After=local-fs.target

[Service]
ExecStart=-/usr/lib/inetutils/telnetd
StandardInput=socket
