# source/xapps-extra/plymouth/rc.zz-plymouth
#!/bin/bash

# (c) 2009 bouleetbil <bouleetbil@frogdev.info>
# rc.zz-plymouth for Frugalware
# distributed under GPL License

# chkconfig: 2345 99 01
# description: Shuts down the splash screen at the end of the boot
source /lib/initscripts/functions
actions=(start stop)

rc_start()
{
	/bin/plymouth --quit
}

rc_stop()
{
	:
}

rc_exec $1
