# source/gnome-extra/system-tools-backends/rc.system-tools-backends
#!/bin/bash

# (c) 2006 Alex Smith <alex.extreme2@gmail.com>
# rc.system-tools-backends for Frugalware
# distributed under GPL License

# chkconfig: 2345 99 01
# description: Backends for Gnome System Tools

source /lib/initscripts/functions
TEXTDOMAIN=system-tools-backends
TEXTDOMAINDIR=/lib/initscripts/messages

actions=(restart start stop)
daemon=$"System Tools Backends"
pid="pidof /usr/bin/system-tools-backends 2> /dev/null"

rc_start()
{
	start_msg
	if [[ "$(eval ${pid})" == "" ]]; then
		/usr/bin/system-tools-backends
		ok $?
	else
		ok 999
	fi
}

rc_stop()
{
	stop_msg
	if [ "$(eval ${pid})" != "" ]; then
		kill $(eval ${pid})
		ok $?
	else
		ok 999
	fi
}

rc_exec $1
