# source/multimedia-extra/lastfmsubmitd/lastfmsubmitd.service
[Unit]
Description=LastFM Submit daemon
After=network.target

[Service]
ExecStart=/usr/bin/lastfmsubmitd
PIDFile=/var/run/lastfm/lastfmsubmitd.pid
Type=forking

[Install]
WantedBy=multi-user.target
