# source/network-extra/openntpd/openntpd.service
[Unit]
Description=OpenNTP Daemon
After=network.target

[Service]
Type=forking
EnvironmentFile=/etc/sysconfig/openntpd
ExecStart=/usr/sbin/ntpd $PARAMS
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
