# source/network-extra/ipw3945d/rc.ipw3945d
#!/bin/bash

# rc.ipw3945 for Frugalware
# distributed under GPL License

# chkconfig: 2345 06 99
# description: ipw3945 daemon

source /lib/initscripts/functions
TEXTDOMAIN=ipw3945d
TEXTDOMAINDIR=/lib/initscripts/messages

actions=(restart start stop)
daemon=$"ipw3945"
PID=`pidof -o %PPID /sbin/ipw3945d`

rc_start()
{
	start_msg
	/sbin/modprobe -C /etc/modprobe.d/ipw3945-add ipw3945
	ok $?
}

rc_stop()
{
	stop_msg
	/sbin/modprobe -C /etc/modprobe.d/ipw3945-rem -r ipw3945
	ok $?
}

rc_exec $1
