# source/network-extra/ddclient/ddclient.service
[Unit]
Description=Dynamic Dns Update Client
After=network.target

[Service]
ExecStart=/usr/sbin/ddclient
Type=forking

[Install]
WantedBy=multi-user.target
