# source/base/net-tools/rc.hostname
#!/bin/sh

# (c) 2003 Vajna Miklos <mamajom@axelero.hu>
# rc.hostname for FrugalWare
# distributed under GPL License

. /etc/rc.d/rc.functions

# chkconfig: 2345 02 98
# description: hostname - set the system's host name

if [ "$1" = "start" ]; then
	# create missing /etc/HOSTNAME if needed
	if ! [ -r /etc/HOSTNAME ]; then
		echo "frugalware.local" >/etc/HOSTNAME
	fi
	hostname=`cat /etc/HOSTNAME | cut -f1 -d .`
	start "$sethostname ($hostname)"
	/bin/hostname "$hostname"
	ok $?
fi
