# source/base/acpid/rc.acpid
#!/bin/sh

# (c) 2003 Vajna Miklos <mamajom@axelero.hu>
# rc.acpid for FrugalWare
# distributed under GPL License

. /etc/rc.d/rc.functions

# chkconfig: 2345 44 56
# description: Listen and dispatch ACPI events from the kernel

[ -d /proc/acpi ] || exit 1

if [ "$1" = "stop" ]; then
	stop "$stopacpi"
	killall acpid
	ok $?
elif [ "$1" = "restart" ]; then
	"$0" stop
	sleep 1
	"$0" start
else # start
	start "$startacpi"
	/usr/sbin/acpid
	ok $?
fi
