# source/base/systemd/rc.local
#!/bin/sh

# rc.local for FrugalWare

#. /etc/rc.d/rc.functions

# chkconfig: 2345 99 01
# description: Local system initialization script

# Put any local setup commands in here:
