# source/network-extra/ntp/ntpd.service
[Unit]
Description=Network Time Service

[Service]
ExecStart=/usr/bin/ntpd -g -c /etc/ntp.org.conf

[Install]
WantedBy=time-sync.target
