# extra/source/network/munin/munin.cron
#
# cron-jobs for munin
#

MAILTO=root

*/5 * * * *     if [ -x /usr/bin/munin-cron ]; then /usr/bin/munin-cron; fi
