# source/network/openssh/sshd-keygen.service
[Unit]
Description=SSH server keys generation.
Before=sshd.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/sshd-keygen
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
